# deploy/systemd/edge-runtime.service
[Unit]
Description=Multi-Persona Voice Agent - Edge Runtime (Raspberry Pi)
Documentation=https://github.com/your-org/multi-persona-voice-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=audio
WorkingDirectory=/home/pi/agent

# Restart policy
Restart=always
RestartSec=10

# Python environment
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/pi/agent"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=edge-runtime

# Run edge runtime
ExecStart=/usr/bin/python3 -m edge.main

# Security/Resource limits
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false  # Need access to home for config
ProtectSystem=strict
ReadWritePaths=/home/pi/agent/data

# Hardware access (needed for camera, audio, Bluetooth)
DevicePolicy=auto
DeviceAllow=/dev/video*
DeviceAllow=/dev/snd/*
DeviceAllow=/dev/dri/*  # For GPU acceleration if available

# Capabilities for hardware access
AmbientCapabilities=
CapabilityBoundingSet=~CAP_SYS_MODULE

[Install]
WantedBy=multi-user.target
